# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ─── Tomato Novel Downloader ───
# Dùng nguyên entrypoint của image (/app/entrypoint.sh) — nó sẽ tự append --server
# + đọc TOMATO_DATA_DIR / TOMATO_WEB_PASSWORD từ env.
[program:tomato]
command=/backend/tomato-init.sh
directory=/app
priority=10
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── Express backend ───
[program:backend]
command=node /backend/src/index.js
directory=/backend
priority=20
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
